# deploy/railway/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# 1. LiveSync bridge: keeps /data/vault and CouchDB in lockstep.
#    Starts before the watcher so the FS is populated.
[program:livesync-bridge]
command=python -m livesync_bridge.cli run
autorestart=true
priority=10
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# 2. Ingestion watcher: reindexes any markdown change into Qdrant.
[program:watcher]
command=python -m ingestion.cli watch
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# 3. MCP server (Phases 2-5): SSE endpoint for Claude.ai + MeetGeek webhook.
[program:mcp-server]
command=python -m mcp_server.server
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
